# templates/network/rc.local.erb
#!/bin/sh -e
#
# Deployment network rc.local for cse480 openstack nodes
# WARNING: this file has been automatically setup by Puppet
# Please make changes there and rerun setup, not here, as they will be overwritten....
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

ifup br-ex
exit 0
